# deploy/alex-backup.service
[Unit]
Description=ALEX NEXUS OS scheduled SQLite backup
After=local-fs.target

[Service]
Type=oneshot
User=vanhkhuc
WorkingDirectory=/opt/alex/AlexRoom-0.2.0-hardware-rc

UMask=0077
Nice=10

ExecStart=/opt/alex/AlexRoom-0.2.0-hardware-rc/.venv/bin/python /opt/alex/AlexRoom-0.2.0-hardware-rc/alex_scheduled_backup.py --database /var/lib/alex/alex.db --backup-dir /var/lib/alex/backups --keep 14

